# scripts/living-captain-status.service
[Unit]
Description=Monad Living Captain status API (read-only, loopback only)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cgl
Group=cgl
WorkingDirectory=/home/cgl/dev/monad
ExecStart=/usr/bin/python3 /home/cgl/dev/monad/tools/living-captain/status_server.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
