# scripts/monad-lan-web.service
[Unit]
Description=Monad LAN-only static web server (web-lan/)
After=network.target

[Service]
Type=simple
User=cgl
Group=cgl
WorkingDirectory=/home/cgl/dev/monad
# Serves web-lan/ (Bridge Station, Bridge Station 2.0, FleetCore Control
# Center, Fleet Motion, Periscope, Radio Console -- see docs/deployment.md's
# Bridge Station 2.0 section for why this is a separate directory and a
# separate process from Caddy's public /var/www/monad root: genuine
# LAN-only isolation needs a port rock64 never forwards, not a path under
# the public root). --bind 0.0.0.0 so it's reachable from any machine on
# Granite's LAN at 192.168.0.100:8090, not just from Granite itself.
ExecStart=/usr/bin/python3 -m http.server 8090 --bind 0.0.0.0 --directory /home/cgl/dev/monad/web-lan
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
