# scripts/rich-voice.service
[Unit]
Description=Monad cache-first Gemini rich voice API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cgl
Group=cgl
WorkingDirectory=/home/cgl/dev/monad
EnvironmentFile=-/home/cgl/.config/monad/gemini.env
Environment=MONAD_VOICE_DAILY_USD=0.10
Environment=MONAD_VOICE_DAILY_SECONDS=300
ExecStart=/usr/bin/python3 /home/cgl/dev/monad/tools/voice-engine/server.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
